# deploy/systemd/video-auto-cut-api.service
[Unit]
Description=Video Auto Cut FastAPI Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=videoautocut
Group=videoautocut
WorkingDirectory=/opt/video_auto_cut
EnvironmentFile=/opt/video_auto_cut/.env
Environment=PYTHONUNBUFFERED=1
Environment=API_HOST=127.0.0.1
Environment=API_PORT=8000
ExecStart=/opt/video_auto_cut/.venv/bin/python -m uvicorn web_api.app:app --host ${API_HOST} --port ${API_PORT} --log-level info --access-log
Restart=always
RestartSec=2
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
